# src/plugin-qt/xsettings/misc/org.deepin.dde.XSettings1.service
[Unit]
Description=org.deepin.dde.XSettings1
Requisite=dde-session-initialized.target dde-session-pre.target
RefuseManualStart=no
RefuseManualStop=no
CollectMode=inactive-or-failed

PartOf=dde-session-pre.target
Before=dde-session-pre.target

Before=dde-session@x11.service

[Service]
Type=dbus
BusName=org.deepin.dde.XSettings1
ExecStart=/usr/bin/deepin-service-manager -n org.deepin.dde.XSettings1
Restart=on-failure
RestartSec=3
Slice=session.slice
NotifyAccess=main
