# src/plugin-qt/xsettings/misc/org.deepin.dde.XSettings1.service
[Unit]
Description=org.deepin.dde.XSettings1
RefuseManualStart=no
RefuseManualStop=no
CollectMode=inactive-or-failed
ConditionEnvironment=XDG_SESSION_TYPE=x11

PartOf=dde-session-pre.target
Before=dde-session-pre.target
After=treeland-sd.service

[Service]
Type=dbus
BusName=org.deepin.dde.XSettings1
ExecStartPre=-/usr/libexec/dde-thp-disable
ExecStart=/usr/bin/deepin-service-manager -n org.deepin.dde.XSettings1
Restart=on-failure
RestartSec=3
Slice=session.slice
NotifyAccess=main
